# ansible/roles/rust_cleanup/templates/rust-cleanup.service.j2
[Unit]
Description=Clean stale Rust target directories under lody worktrees

[Service]
Type=oneshot
Nice=10
IOSchedulingClass=idle
ExecStart=/bin/sh -lc '[ -d "%h/.lody/repos" ] || exit 0; find "%h/.lody/repos" -mindepth 3 -maxdepth 3 -type d -path "*/worktrees/*" -print0 | xargs -0 -r -n 1 "{{ rust_cleanup_entrypoint }}"'
